# payload/Library/LaunchAgents/com.github.uptimenag.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.github.uptimenag</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/uptime_nag/check_uptime.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>ThrottleInterval</key>
    <integer>0</integer>
    <key>KeepAlive</key>
    <false/>
    <key>StandardErrorPath</key>
    <string>/ApplicationSupport/UptimeNag/out_log.log</string>
    <key>StandardOutPath</key>
    <string>/ApplicationSupport/UptimeNag/err_log.log</string>
</dict>
</plist>